# src/init.d/sysadm-rest
#!/sbin/openrc-run
# Copyright (c) 2016 Ken Moore <ken@ixsystems.com>
# Released under the 2-clause BSD license

name="SysAdm Server (REST)"
pidfile=/var/run/sysadm-rest.pid
prefix=/usr/local/
command=${prefix}bin/sysadm-binary
command_args="--rest"
supervisor=supervise-daemon

depend() { 
	after bootmisc
	keyword -shutdown
}

ssl_keygen()
{
	if [ ! -d "${prefix}/etc/sysadm" ] ; then
    		mkdir -p ${prefix}/etc/sysadm
	fi
	openssl req -x509 -nodes -newkey rsa:2048 \
		-keyout ${prefix}/etc/sysadm/restserver.key \
		-out ${prefix}/etc/sysadm/restserver.crt -days 1024 \
		-subj "/C=US/ST=MY/L=NULL/O=SysAdm/OU=SysAdm/CN=SysAdm/emailAddress=none@example.org" 2>/dev/null
}

start_pre()
{
	ssl_keygen "$1"
}

stop()
{
        t=$(mktemp)
        for pid in $(cat $pidfile) ; do
          echo $pid >$t
          pidfile=$t default_stop
        done
        rm -f $t
}
