# src/port-files/pkg-plist
bin/sysadm-server
etc/rc.d/sysadm
etc/rc.d/sysadm-rest
