# src/port-files/pkg-plist
bin/sysadm-server
etc/rc.d/sysadm-restserver
etc/rc.d/sysadm-websocket
