# port-files/pkg-plist
bin/sysadm-binary
bin/sysadm-server
bin/sysadm-bridge
bin/sysadm-bridge-start
etc/sysadm.conf.dist
etc/rc.d/sysadm
etc/rc.d/sysadm-bridge
etc/rc.d/sysadm-rest
