# port-files/files/sysadm-bridge.in
#!/sbin/openrc-run
# Copyright (c) 2016 Ken Moore <ken@ixsystems.com>
# Released under the 2-clause BSD license

command="%%PREFIX%%/bin/sysadm-bridge"
command_args="&"
pid=/var/run/sysadm-bridge.pid
name="SysAdm Bridge"

depend() { 
	need syscons net
	after bootmisc
	keyword -shutdown
}
