# port-files/files/sysadm-bridge.in
#!/sbin/openrc-run
# Copyright (c) 2016 Ken Moore <ken@ixsystems.com>
# Released under the 2-clause BSD license

name="SysAdm Server (Bridge)"

depend() { 
	after bootmisc
	keyword -shutdown
}

start()
{
	ebegin "Starting $name"
	start-stop-daemon --start -m -b -p /var/run/sysadm-bridge.pid \
		%%PREFIX%%/bin/sysadm-bridge
	eend $?
}

stop()
{
	ebegin "Stopping $name"
	start-stop-daemon --stop --exec %%PREFIX%%/bin/sysadm-bridge \
		-p /var/run/sysadm-bridge.pid
	eend $?
}

ssl_keygen()
{
	if [ ! -d "%%PREFIX%%/etc/sysadm" ] ; then
    		mkdir -p %%PREFIX%%/etc/sysadm
	fi
	if [ ! -e "%%PREFIX%%/etc/sysadm/ws_bridge.key" ] ; then
		openssl req -x509 -nodes -newkey rsa:2048 \
		-keyout %%PREFIX%%/etc/sysadm/ws_bridge.key \
		-out %%PREFIX%%/etc/sysadm/ws_bridge.crt -days 102400 \
		-subj "/C=US/ST=MY/L=NULL/O=SysAdm/OU=SysAdm/CN=SysAdm/emailAddress=none@example.org" 2>/dev/null
	fi
}

start_pre()
{
	ssl_keygen "$1"
}
